# sysd/pishop.crontab
0 5 0 0 0 /sbin/shutdown -r now
@reboot rm /home/pi/.cache/chromium/Default/Cache/*
